# deploy/lite/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

; =============================================================================
; Redis Server (internal)
; =============================================================================
[program:redis]
command=/usr/bin/redis-server --bind 0.0.0.0 --port 6379 --dir /var/lib/redis --appendonly yes --appendfsync everysec
autostart=true
autorestart=true
priority=5
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; Xvfb - Virtual X11 Display (for Playwright browsers)
; =============================================================================
[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; PulseAudio - Audio server for TTS playback
; =============================================================================
[program:pulseaudio]
command=pulseaudio --system --disallow-exit --disallow-module-loading=false --log-level=notice
autostart=true
autorestart=true
priority=12
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; PulseAudio Sink Setup (runs once after PulseAudio is ready)
; =============================================================================
[program:pulseaudio-setup]
command=/usr/local/bin/setup-pulseaudio-sinks.sh
autostart=true
autorestart=false
startsecs=0
priority=13
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; x11vnc — VNC server sharing the Xvfb display for browser sessions
; =============================================================================
[program:x11vnc]
command=x11vnc -display :99 -forever -nopw -shared -rfbport 5900 -q
autostart=true
autorestart=true
priority=14
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; websockify — bridges VNC (port 5900) to WebSocket (port 6080) for noVNC
; =============================================================================
[program:websockify]
command=websockify --web /usr/share/novnc 6080 localhost:5900
autostart=true
autorestart=true
priority=15
startsecs=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

; =============================================================================
; Admin API (:8057) — user and API token management
; =============================================================================
[program:admin-api]
command=uvicorn app.main:app --host 0.0.0.0 --port 8057 --log-level %(ENV_LOG_LEVEL)s
directory=/app/admin-api
autostart=true
autorestart=true
priority=30
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/admin-models:/app/meeting-api:/app/schema-sync",DATABASE_URL="%(ENV_DATABASE_URL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",ADMIN_API_TOKEN="%(ENV_ADMIN_API_TOKEN)s",INTERNAL_API_SECRET="lite-internal-secret",REDIS_URL="%(ENV_REDIS_URL)s",TRANSCRIPTION_SERVICE_URL="%(ENV_TRANSCRIPTION_SERVICE_URL)s",TRANSCRIPTION_SERVICE_TOKEN="%(ENV_TRANSCRIPTION_SERVICE_TOKEN)s"

; =============================================================================
; Runtime API (:8090) — container lifecycle with process backend
; =============================================================================
[program:runtime-api]
command=uvicorn runtime_api.main:app --host 0.0.0.0 --port 8090 --log-level %(ENV_LOG_LEVEL)s
directory=/app/runtime-api
autostart=true
autorestart=true
priority=35
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",ORCHESTRATOR_BACKEND="process",PORT="8090",REDIS_URL="%(ENV_REDIS_URL)s",BOT_SCRIPT_PATH="/app/vexa-bot/dist/docker.js",BOT_WORKING_DIR="/app/vexa-bot",DISPLAY=":99",ALLOW_PRIVATE_CALLBACKS="true"

; =============================================================================
; Meeting API (:8080) — bot orchestration + transcription pipeline
; =============================================================================
[program:meeting-api]
command=uvicorn meeting_api.main:app --host 0.0.0.0 --port 8080 --log-level %(ENV_LOG_LEVEL)s
directory=/app/meeting-api
autostart=true
autorestart=true
priority=40
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/admin-models:/app/schema-sync",REDIS_URL="%(ENV_REDIS_URL)s",DATABASE_URL="%(ENV_DATABASE_URL)s",DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s",RUNTIME_API_URL="http://localhost:8090",MEETING_API_URL="http://localhost:8080",BOT_CALLBACK_BASE_URL="http://localhost:8080",ADMIN_API_URL="http://localhost:8057",ADMIN_TOKEN="%(ENV_ADMIN_API_TOKEN)s",TRANSCRIPTION_SERVICE_URL="%(ENV_TRANSCRIPTION_SERVICE_URL)s",TRANSCRIPTION_SERVICE_TOKEN="%(ENV_TRANSCRIPTION_SERVICE_TOKEN)s",TTS_SERVICE_URL="http://localhost:8059",STORAGE_BACKEND="%(ENV_STORAGE_BACKEND)s",LOCAL_STORAGE_DIR="%(ENV_LOCAL_STORAGE_DIR)s",LOCAL_STORAGE_FSYNC="%(ENV_LOCAL_STORAGE_FSYNC)s",MINIO_ENDPOINT="%(ENV_MINIO_ENDPOINT)s",MINIO_ACCESS_KEY="%(ENV_MINIO_ACCESS_KEY)s",MINIO_SECRET_KEY="%(ENV_MINIO_SECRET_KEY)s",MINIO_BUCKET="%(ENV_MINIO_BUCKET)s",MINIO_SECURE="%(ENV_MINIO_SECURE)s",AWS_REGION="%(ENV_AWS_REGION)s",AWS_ACCESS_KEY_ID="%(ENV_AWS_ACCESS_KEY_ID)s",AWS_SECRET_ACCESS_KEY="%(ENV_AWS_SECRET_ACCESS_KEY)s",S3_BUCKET="%(ENV_S3_BUCKET)s",S3_ENDPOINT="%(ENV_S3_ENDPOINT)s",S3_SECURE="%(ENV_S3_SECURE)s"

; =============================================================================
; Agent API (:8100) — AI agent chat runtime
; =============================================================================
[program:agent-api]
command=uvicorn agent_api.main:app --host 0.0.0.0 --port 8100 --log-level %(ENV_LOG_LEVEL)s
directory=/app/agent-api
autostart=true
autorestart=true
priority=45
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",REDIS_URL="%(ENV_REDIS_URL)s",RUNTIME_API_URL="http://localhost:8090",MEETING_API_URL="http://localhost:8080"

; =============================================================================
; TTS Service (:8059) — text-to-speech via OpenAI API
; =============================================================================
[program:tts-service]
command=uvicorn main:app --host 0.0.0.0 --port 8059 --log-level %(ENV_LOG_LEVEL)s
directory=/app/tts-service
autostart=true
autorestart=true
priority=45
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s"

; =============================================================================
; API Gateway (:8056) — main entry point, routes to all services
; =============================================================================
[program:api-gateway]
command=uvicorn main:app --host 0.0.0.0 --port 8056 --log-level %(ENV_LOG_LEVEL)s
directory=/app/api-gateway
autostart=true
autorestart=true
priority=50
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app/meeting-api:/app/admin-models:/app/schema-sync",ADMIN_API_URL="http://localhost:8057",MEETING_API_URL="http://localhost:8080",TRANSCRIPTION_COLLECTOR_URL="http://localhost:8080",MCP_URL="http://localhost:18888",AGENT_API_URL="http://localhost:8100",REDIS_URL="%(ENV_REDIS_URL)s",INTERNAL_API_SECRET="lite-internal-secret",VNC_HOST="localhost"

; =============================================================================
; MCP Service (:18888) — Model Context Protocol for Claude/Cursor
; =============================================================================
[program:mcp]
command=uvicorn main:app --host 0.0.0.0 --port 18888 --log-level %(ENV_LOG_LEVEL)s
directory=/app/mcp
autostart=true
autorestart=true
priority=55
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",API_GATEWAY_URL="http://localhost:8056",LOG_LEVEL="%(ENV_LOG_LEVEL)s"

; =============================================================================
; Dashboard (:3000) — Next.js web UI
; =============================================================================
[program:dashboard]
command=node server.js
directory=/app/dashboard
autostart=true
autorestart=true
priority=60
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",NEXT_PUBLIC_API_URL="http://localhost:8056",VEXA_API_URL="http://localhost:8056",VEXA_ADMIN_API_KEY="%(ENV_ADMIN_API_TOKEN)s",VEXA_ADMIN_API_URL="http://localhost:8056"

; =============================================================================
; Process Groups
; =============================================================================
[group:vexa-core]
programs=redis,xvfb,pulseaudio,pulseaudio-setup,x11vnc,websockify,admin-api,runtime-api,meeting-api,agent-api,tts-service,api-gateway,mcp,dashboard
priority=999
